# systemd/website-watcher.timer
[Unit]
Description=Run Website Watcher every 2 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=2h
Persistent=true

[Install]
WantedBy=timers.target
